# pvs-labs/lab-rip/net/r4/etc/rc.local
#!/bin/sh -e

iptables -F 
iptables -F -t nat
iptables -X

iptables -t nat -A POSTROUTING -o eth2 -j MASQUERADE

exit 0
